# debian/vyatta-cfg-system.postinst.in
#!/bin/bash

prefix=@prefix@
sysconfdir=@sysconfdir@

# remove init of daemons that we manage
update-rc.d -f ntp remove >/dev/null

