# debian/vyatta-cfg-system.postrm
#!/bin/bash

if [ "$1" = "purge" ]; then
    sed -i -e 'g/^password/d' /etc/pam.d/password
    update-rc.d vyatta-config-reboot-params remove
fi

# Local Variables:
# mode: shell-script
# sh-indentation: 4
# End:
