# etc/profile.d/vyatta-login.sh
/opt/vyatta/sbin/vyatta-first-login-passwd.sh
